# debian/postinst
# Disable dhcpcd from launching wpa_supplicant
rm /usr/lib/dhcpcd/dhcpcd-hooks/10-wpa_supplicant

# Update wpa_supplicant to launch with dbus support
sed -i 's%ExecStart=/sbin/wpa_supplicant -c%ExecStart=/sbin/wpa_supplicant -u -c%' /usr/lib/systemd/system/wpa_supplicant@.service
systemctl disable wpa_supplicant
systemctl enable wpa_supplicant@wlan0

useradd -m bartender
chown -R bartender:bartender /usr/share/bartender
chown -R bartender:bartender /home/bartender
usermod -a -G i2c bartender
usermod -a -G tty bartender
#usermod -a -G gpio bartender