# debian/postrm
# Update wpa_supplicant to launch without dbus support (the default)
sed -i 's%ExecStart=/sbin/wpa_supplicant -u -c%ExecStart=/sbin/wpa_supplicant -c%' /usr/lib/systemd/system/wpa_supplicant@.service
systemctl disable wpa_supplicant@wlan0
systemctl enable wpa_supplicant
